# systemd/botzoom-backend.service
[Unit]
Description=BotZoom FastAPI Backend Service
After=network.target

[Service]
Type=simple
User=cak-seno
Group=cak-seno
WorkingDirectory=/home/cak-seno/botzoom/backend
Environment="PATH=/home/cak-seno/botzoom/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DISPLAY=:99"
Environment="DBUS_SESSION_BUS_ADDRESS=/dev/null"
Environment="PULSE_SERVER=/run/user/1000/pulse/native"
Environment="XDG_RUNTIME_DIR=/run/user/1000"

# Load environment variables from .env.production
EnvironmentFile=/home/cak-seno/botzoom/backend/.env.production

# Start command (1 worker to prevent OOM - Whisper model uses ~1.5GB per worker)
ExecStart=/home/cak-seno/botzoom/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/botzoom/backend.log
StandardError=append:/var/log/botzoom/backend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
